# ops/systemd/tracker-api.service
[Unit]
Description=KSR Tracker API
After=network.target

[Service]
Type=simple
User=tracker
Group=tracker
WorkingDirectory=/opt/ksr-tracker
EnvironmentFile=/opt/ksr-tracker/.env
ExecStart=/opt/ksr-tracker/.venv/bin/uvicorn server.app.main:app --host 127.0.0.1 --port 8000 --proxy-headers --forwarded-allow-ips=*
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
